# deployment/init.d/simple-ntpd
#!/bin/bash
#
# simple-ntpd    Simple NTP Daemon init script
#
# chkconfig: 2345 20 80
# description: Simple NTP Daemon provides NTP time synchronization
# processname: simple-ntpd
# pidfile: /var/run/simple-ntpd/simple-ntpd.pid
# config: /etc/simple-ntpd/simple-ntpd.conf
#
# Author: SimpleDaemons
# License: Apache 2.0

### BEGIN INIT INFO
# Provides: simple-ntpd
# Required-Start: $network $remote_fs $syslog
# Required-Stop: $network $remote_fs $syslog
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Short-Description: Simple NTP Daemon
# Description: Simple NTP Daemon provides NTP time synchronization
### END INIT INFO

# Source function library
. /etc/rc.d/init.d/functions

# Source networking configuration
. /etc/sysconfig/network

# Check that networking is up
[ "${NETWORKING}" = "no" ] && exit 6

# NTP daemon configuration
NTPD_NAME="simple-ntpd"
NTPD_BIN="/usr/sbin/simple-ntpd"
NTPD_CONF="/etc/simple-ntpd/simple-ntpd.conf"
NTPD_PID="/var/run/simple-ntpd/simple-ntpd.pid"
NTPD_LOCK="/var/lock/subsys/simple-ntpd"
NTPD_LOG="/var/log/simple-ntpd/simple-ntpd.log"
NTPD_USER="ntp"
NTPD_GROUP="ntp"
NTPD_ARGS="-c ${NTPD_CONF}"

# Check if binary exists
[ -x ${NTPD_BIN} ] || exit 5

# Check if config exists
[ -f ${NTPD_CONF} ] || exit 6

# Create necessary directories
create_dirs() {
    mkdir -p /var/run/simple-ntpd
    mkdir -p /var/log/simple-ntpd
    mkdir -p /var/lib/simple-ntpd
    chown ${NTPD_USER}:${NTPD_GROUP} /var/run/simple-ntpd
    chown ${NTPD_USER}:${NTPD_GROUP} /var/log/simple-ntpd
    chown ${NTPD_USER}:${NTPD_GROUP} /var/lib/simple-ntpd
    chmod 755 /var/run/simple-ntpd
    chmod 755 /var/log/simple-ntpd
    chmod 755 /var/lib/simple-ntpd
}

# Start the daemon
start() {
    echo -n "Starting ${NTPD_NAME}: "
    
    # Create necessary directories
    create_dirs
    
    # Check if already running
    if [ -f ${NTPD_PID} ]; then
        if kill -0 $(cat ${NTPD_PID}) 2>/dev/null; then
            echo "already running"
            return 1
        else
            rm -f ${NTPD_PID}
        fi
    fi
    
    # Start the daemon
    daemon --user ${NTPD_USER} --pidfile ${NTPD_PID} ${NTPD_BIN} ${NTPD_ARGS}
    RETVAL=$?
    
    if [ $RETVAL -eq 0 ]; then
        echo "OK"
        touch ${NTPD_LOCK}
    else
        echo "FAILED"
    fi
    
    return $RETVAL
}

# Stop the daemon
stop() {
    echo -n "Stopping ${NTPD_NAME}: "
    
    if [ -f ${NTPD_PID} ]; then
        killproc -p ${NTPD_PID} ${NTPD_BIN}
        RETVAL=$?
        if [ $RETVAL -eq 0 ]; then
            rm -f ${NTPD_PID}
            rm -f ${NTPD_LOCK}
            echo "OK"
        else
            echo "FAILED"
        fi
    else
        echo "not running"
        RETVAL=0
    fi
    
    return $RETVAL
}

# Restart the daemon
restart() {
    stop
    sleep 2
    start
}

# Reload configuration
reload() {
    echo -n "Reloading ${NTPD_NAME}: "
    
    if [ -f ${NTPD_PID} ]; then
        if kill -HUP $(cat ${NTPD_PID}) 2>/dev/null; then
            echo "OK"
            return 0
        else
            echo "FAILED"
            return 1
        fi
    else
        echo "not running"
        return 1
    fi
}

# Check status
status() {
    if [ -f ${NTPD_PID} ]; then
        if kill -0 $(cat ${NTPD_PID}) 2>/dev/null; then
            echo "${NTPD_NAME} is running (pid: $(cat ${NTPD_PID}))"
            return 0
        else
            echo "${NTPD_NAME} is not running (stale pid file)"
            return 1
        fi
    else
        echo "${NTPD_NAME} is not running"
        return 3
    fi
}

# Check configuration
configtest() {
    echo -n "Testing ${NTPD_NAME} configuration: "
    
    if [ -f ${NTPD_CONF} ]; then
        if ${NTPD_BIN} -t -c ${NTPD_CONF} >/dev/null 2>&1; then
            echo "OK"
            return 0
        else
            echo "FAILED"
            return 1
        fi
    else
        echo "configuration file not found"
        return 1
    fi
}

# Show usage
usage() {
    echo "Usage: $0 {start|stop|restart|reload|status|configtest}"
    echo "  start       - Start the NTP daemon"
    echo "  stop        - Stop the NTP daemon"
    echo "  restart     - Restart the NTP daemon"
    echo "  reload      - Reload configuration"
    echo "  status      - Show daemon status"
    echo "  configtest  - Test configuration file"
    exit 1
}

# Main script logic
case "$1" in
    start)
        start
        ;;
    stop)
        stop
        ;;
    restart)
        restart
        ;;
    reload)
        reload
        ;;
    status)
        status
        ;;
    configtest)
        configtest
        ;;
    *)
        usage
        ;;
esac

exit $?
