# packaging/linux/deb/prerm
#!/bin/bash
# Pre-removal script for simple-ntpd

set -e

PROJECT_NAME="simple-ntpd"

if [ "$1" = "remove" ] || [ "$1" = "deconfigure" ]; then
    # Stop service before removal
    systemctl stop "$simple-ntpd" 2>/dev/null || true
    systemctl disable "$simple-ntpd" 2>/dev/null || true
fi

exit 0

